# deploy/icsdpanel.service
[Unit]
Description=ICSD Panel - Linux server management panel
Documentation=https://github.com/ICSD/ICSDPanel
After=network.target

[Service]
Type=simple
User=icsdpanel
WorkingDirectory=/opt/icsdpanel/backend
EnvironmentFile=/opt/icsdpanel/backend/.env
ExecStart=/opt/icsdpanel/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8088
Restart=on-failure
RestartSec=5
# Hardening
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
